# docker/cron/crontab
* * * * * www-data /usr/local/bin/php /var/www/railt.org/artisan schedule:run >> /var/log/cron.log 2>&1
